# dockerconf/supervisord.conf
[supervisord]
nodaemon=true

[program:dws]
command=npm start
autostart=true
autorestart=true
