# dockerconf/supervisord.conf
[supervisord]
nodaemon=true

[program:node]
directory=/app
command=node boot.js
autostart=true
autorestart=true
